# crontab
* * * * * cd /home/nthansen/linux-hanse/side-projects/game-scraper/ && .venv/bin/python3 ./webscraper.py